# server-image/supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
command=/usr/sbin/tor -f /etc/tor/torrc

[program:privoxy]
command=/usr/sbin/privoxy --no-daemon /etc/privoxy/config
